# templates/etc/profile.d/etcd.sh
ETCD_INITIAL_CLUSTER="ubuntu-bionic-1=http://10.128.0.2:2380,ubuntu-bionic-2=http://10.128.0.3:2380,ubuntu-bionic-3=http://10.128.0.4:2380"
ETCD_INITIAL_CLUSTER_STATE=new