# templates/etc/profile.d/etcd.sh
ETCD_INITIAL_CLUSTER="infra0=http://10.128.0.2:2380,infra1=http://10.128.0.3:2380,infra2=http://10.128.0.4:2380"
ETCD_INITIAL_CLUSTER_STATE=new